# meta-senic/recipes-hub/hub-configuration/files/supervisor_senic_hub.conf
[program:senic_hub]
command = /usr/bin/pserve /etc/senic_hub.ini
directory=/data/senic-hub
autostart=true
autorestart=true
stdout_logfile=/data/senic-hub/logs/senic_hub.log
redirect_stderr=true
stdout_logfile_backup=1
stdout_capture_maxbytes=50MB
stopsignal=QUIT
user=root
environment=LC_ALL=en_US.UTF-8, LANG=en_US.UTF-8
