# meta-senic/recipes-connectivity/enable-bt-module-systemd/files/enable-bt-module.service
[Unit]
Description=Enable bt module
[Service]
Type=simple
ExecStart=/usr/bin/enable-bt-module
[Install]
WantedBy=multi-user.target bluetooth.service
